# config/etc/cron.d/api_command-cron
0 * * * * /usr/local/bin/api_command.sh
